# deploy/verigo-retention.timer
[Unit]
Description=Daily Verigo data retention cleanup

[Timer]
OnCalendar=*-*-* 04:10:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
