# deploy/verigo-backup.timer
[Unit]
Description=Daily Verigo backup

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
